# infra/systemd/ark-portal-inventory.service
[Unit]
Description=Ark Portal — inventory (inventory.arkinstitutebc.com)
After=network.target

[Service]
Type=simple
User=ark
Group=ark
WorkingDirectory=/opt/ark-portals/repo/apps/inventory
ExecStart=/home/ark/.bun/bin/bun dist/server/index.mjs
EnvironmentFile=/opt/ark-portals/repo/apps/inventory/.env.production
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
