# infra/systemd/ark-portal-billing.service
[Unit]
Description=Ark Portal — billing (billing.arkinstitutebc.com)
After=network.target

[Service]
Type=simple
User=ark
Group=ark
WorkingDirectory=/opt/ark-portals/repo/apps/billing
ExecStart=/home/ark/.bun/bin/bun dist/server/index.mjs
EnvironmentFile=/opt/ark-portals/repo/apps/billing/.env.production
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
