# infra/systemd/ark-portal-training.service
[Unit]
Description=Ark Portal — training (training.arkinstitutebc.com)
After=network.target

[Service]
Type=simple
User=ark
Group=ark
WorkingDirectory=/opt/ark-portals/repo/apps/training
ExecStart=/home/ark/.bun/bin/bun dist/server/index.mjs
EnvironmentFile=/opt/ark-portals/repo/apps/training/.env.production
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
